# server/templates/supervisor.ini
[program:{{ app.name }}]
command=/usr/local/bin/uwsgi
  --socket 127.0.0.1:5000
  --logto /home/app/log/application.log
  --home /home/app/env 
  --pythonpath /home/app/src
  --wsgi-file /home/app/src/{{ app.wsgi_file }}
  --callable {{ app.wsgi_var_name }}
  --max-requests 1000
  --master
  --processes {{ app.workers }}
  --chmod 
user=app
directory=/home/app/src
autostart=true
autorestart=true
